# etc/systemd/system/serial-getty@ttySAC2.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin root --keep-baud 115200 %I $TERM
