# upstart/waggle-epoch.conf
start on (local-filesystems)
stop on runlevel[016]

respawn
respawn limit unlimited
post-stop exec sleep 5

chdir /usr/lib/waggle/core/

exec /usr/lib/waggle/core/scripts/waggle_epoch.sh
