# upstart/waggle-heartbeat.conf
start on (local-filesystems)
stop on runlevel[016]
respawn
respawn limit unlimited
post-stop exec sleep 5
chdir /usr/lib/waggle/core/scripts
exec /usr/lib/waggle/core/scripts/heartbeat.sh

